# green.conf
[program:green_card]

directory=home/ubuntu/green_card_django
command=